# deploy/systemd/triparr-web.service
[Unit]
Description=Triparr Web Application
After=network.target

[Service]
Type=simple
User=triparr
WorkingDirectory=/opt/triparr
ExecStart=/usr/bin/env npm start
Restart=always
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
